# deployments/studio/mlx_server.service
[Unit]
Description=Monolathe MLX Inference Server
After=network.target

[Service]
Type=simple
User=studio
Group=staff
WorkingDirectory=/Users/studio/monolathe
Environment=PATH=/Users/studio/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONPATH=/Users/studio/monolathe
Environment=MAX_VRAM_USAGE_GB=44
Environment=MAX_CONCURRENT_VIDEO_GENS=2
Environment=MAX_CONCURRENT_IMAGE_GENS=4
Environment=SHARED_STORAGE_PATH=/Volumes/ai_content_shared

ExecStart=/Users/studio/.local/bin/uvicorn src.assetfactory.mlx_server:app \
    --host 0.0.0.0 \
    --port 8080 \
    --workers 1 \
    --loop uvloop

Restart=always
RestartSec=10

# Resource limits
MemoryMax=48G
CPUQuota=800%

[Install]
WantedBy=multi-user.target
